# blfs-scripts/b/bash-profile-7.6/scripts/etc/profile.d/i18n.sh
# This script sets an environment variable necessary for native language 
# support. A full discussion on determining this variable can be found on 
# the LFS Bash Shell Startup Files page.

# Set up i18n variables
export LANG=en_US.UTF-8
